# machine-learning/inference-on-ec2/qwen3-tts-aws-sample/scripts/qwen3-tts.service
[Unit]
Description=Qwen3-TTS Gradio Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/app
Environment="HF_HOME=/opt/huggingface"
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=/usr/bin/python3 /opt/app/server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/gradio-server.log
StandardError=append:/var/log/gradio-server.log

[Install]
WantedBy=multi-user.target
